# systemd/bella-htw.automount
[Unit]
Description=Automount BELLA HTW

[Automount]
Where=/bella/htw

[Install]
WantedBy=multi-user.target
