# systemd/bella-pw.automount
[Unit]
Description=Automount BELLA PW

[Automount]
Where=/bella/pw

[Install]
WantedBy=multi-user.target
